# Combined-Leaderboard/deployment/huggingface/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/ms-vista-supervisord.pid
logfile=/dev/null
logfile_maxbytes=0

[program:api]
command=/usr/local/bin/gunicorn backend.web:app --bind 127.0.0.1:%(ENV_API_PORT)s --workers %(ENV_WEB_CONCURRENCY)s --worker-class gthread --threads %(ENV_GUNICORN_THREADS)s --timeout %(ENV_GUNICORN_TIMEOUT)s
directory=/home/user/app
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -c /home/user/app/deployment/huggingface/nginx.conf -g "daemon off;"
directory=/home/user/app
autostart=true
autorestart=true
startsecs=1
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
